# cron/postgresql-partman-to-cstore-move-data.cron
#!/bin/bash

/usr/bin/partman_to_cstore_move_data --quiet
:
